# storebot.service
[Unit]
Description=Telegram Store Bot
After=network.target

[Service]
Type=simple
User=YOUR_VPS_USERNAME
WorkingDirectory=/home/YOUR_VPS_USERNAME/telegram-store-bot
Environment="PATH=/home/YOUR_VPS_USERNAME/telegram-store-bot/venv/bin"
ExecStart=/home/YOUR_VPS_USERNAME/telegram-store-bot/venv/bin/python storebot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
